# misc/pavd.service
[Unit]
Description=secPaver daenmon

[Service]
ExecStart=/usr/bin/pavd
Restart=on-failure

[Install]
WantedBy=multi-user.target